# docker/build/final/cron/etc/crontab
#crontab running scripts into exist containers
*	*	*	*	*	docker exec -u www-data:www-data php php -v

